# supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=/usr/sbin/cron -f

[program:ghost]
command=node current/index.js
directory=/var/lib/ghost
